# backup/crontab
# Daily backup at 03:00 Asia/Almaty
0 3 * * * /scripts/backup.sh >> /var/log/backup.log 2>&1
